# audio-switch.service
[Unit]
Description=Audio switch
After=network.target
StartLimitIntervalSec=10

[Service]
Type=simple
Restart=always
RestartSec=1
User=stefan
ExecStart=/home/stefan/audio-switch/.venv/bin/python /home/stefan/audio-switch/main.py

[Install]
WantedBy=multi-user.target